# deploy/qc-dashboard.service
[Unit]
Description=QC Incoming Material Inspection — Web Dashboard
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/.venv/bin/python src/web_dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
